# cms.service
[Unit]
Description=CMS Django Application
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target postgresql.service redis-server.service

[Service]
Type=notify
User=cmsuser
Group=www-data
WorkingDirectory=/home/cmsuser/cms_project
Environment="PATH=/home/cmsuser/cms_project/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DJANGO_SETTINGS_MODULE=cms_project.settings"
Environment="PYTHONPATH=/home/cmsuser/cms_project"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONDONTWRITEBYTECODE=1"
Environment="HOME=/home/cmsuser"
RuntimeDirectory=cms
RuntimeDirectoryMode=0750

# Security hardening (safe subset — Python + gunicorn compatible)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/run/cms /home/cmsuser/cms_project/logs /home/cmsuser/cms_project/media /home/cmsuser/cms_project/staticfiles
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Pre-start validation — run with service-manager privileges so ownership
# fixes do not fail before Gunicorn can start.
ExecStartPre=+/bin/mkdir -p /home/cmsuser/cms_project/logs /home/cmsuser/cms_project/media /home/cmsuser/cms_project/staticfiles
ExecStartPre=+/bin/chown -R cmsuser:www-data /home/cmsuser/cms_project/logs /home/cmsuser/cms_project/media /home/cmsuser/cms_project/staticfiles
ExecStartPre=+/bin/chmod -R 755 /home/cmsuser/cms_project/logs /home/cmsuser/cms_project/media /home/cmsuser/cms_project/staticfiles

# Process management
ExecStart=/home/cmsuser/cms_project/venv/bin/python -m gunicorn \
    --config /home/cmsuser/cms_project/gunicorn.conf.py \
    cms_project.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
KillMode=mixed
TimeoutStopSec=30
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cms

# Notify systemd when ready
NotifyAccess=all

[Install]
WantedBy=multi-user.target
